# usage-process/trunk/sbin/move_to_repository.cron
0 6 * * * /soft/usage-process/sbin/move_to_repository.py -c /soft/usage-process/conf/move_to_repository.conf -l info
