# usage-process/sbin/move_to_repository.cron
0 13 * * * /soft/XCI-Usage-Tools/PROD/usage-process/sbin/move_to_repository.py -c /soft/XCI-Usage-Tools/conf/move_to_repository.conf -l info
